# services/motioneye/finish
#!/bin/bash

meyectl stopserver /config/motioneye.conf
