# services/motioneye/finish
#!/bin/bash

meyectl stopserver -c /config/motioneye.conf
