# supervisor/motion.conf
[program:motion]
user=nobody
command=motion
stopwaitsecs=30
