# supervisor/motion.conf
[program:motion]
user=nobody
command=exec motion
stopwaitsecs=30
