# supervisor/motioneye.conf
[program:motioneye]
user=nobody
command=meyectl startserver -b -c /config/motioneye.conf
stopwaitsecs=30
