# supervisor/motioneye.conf
[program:motioneye]
user=nobody
command=bash -c 'sleep 15 && meyectl startserver -b -c /config/motioneye.conf'
stopwaitsecs=30
