# supervisor/motion.conf
[program:motion]
user=nobody
command=/usr/local/bin/dockmotion
stopwaitsecs=30
